# systemd/halo-lemonade.service
[Unit]
Description=Lemonade — Intelligent LLM Gateway
After=network.target halo-llama-server.service
Wants=halo-llama-server.service

[Service]
Type=simple
User=<YOUR_USER>
Group=<YOUR_USER>
WorkingDirectory=/srv/ai/lemonade
ExecStart=/usr/local/bin/node /srv/ai/lemonade/dist/index.js
Environment=LEMONADE_HOME=/srv/ai/lemonade
Environment=PORT=8080
Restart=on-failure
RestartSec=5
ProtectSystem=strict
PrivateTmp=yes
NoNewPrivileges=yes
ProtectHome=read-only
ReadWritePaths=/srv/ai

[Install]
WantedBy=multi-user.target
