# iso/airootfs/etc/systemd/system/halo-autoinstall.service
[Unit]
Description=Halo AI Core — Automated Installer
After=network-online.target
Wants=network-online.target
ConditionKernelCommandLine=halo-autoinstall

[Service]
Type=oneshot
ExecStart=/opt/halo-autoinstall/autoinstall.sh --unattended
StandardOutput=journal+console
StandardError=journal+console
TTYPath=/dev/tty1
TimeoutStartSec=3600

[Install]
WantedBy=multi-user.target
